# crontab
*/5 * * * * root rake submission:generate-pdfs
00 00 * * * root rake db:populate SCALE=large EXTENDED=true
